# supervisor.conf
[supervisord]
nodaemon=true

[inet_http_server]
port=*:9001

[program:apache]
command=/usr/sbin/apache2 -DFOREGROUND

[program:not_cron]
command=/var/www/not_cron.sh
directory=/var/www

